# packages/uun-livecam/uun_livecam-0.6.0-py3-none-any.whl/uun_livecam/data/systemd.service
[Unit]
Description=uunLiveCam Gateway Service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=~
ExecStart=%h/start.sh -l WARNING --systemd

Restart=on-failure
RestartSec=30s

ProtectSystem=strict
ReadWritePaths=%h

[Install]
WantedBy=default.target

